# image_builder/python/uwsgi_sample_with_pyrunit_and_nginx/etc/sv/app/run
#!/bin/bash
exec 2>&1
cd /usr/src/app
exec chpst -u root uwsgi --ini app.ini
